# files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
[program:{{ service_name }}]
command=node /srv/statsd/statsd/stats.js /etc/statsd/localConfig.js
user=statsd
stdout_logfile=/var/log/statsd/application.log
stderr_logfile=/var/log/statsd/error.log
autostart=true
autorestart=true
redirect_stderr=true
